# deploy/plan1.service
[Unit]
Description=plan1 (sillyz.computer peer)
After=network.target

[Service]
User=clownbot
WorkingDirectory=/home/clownbot/plan1
Environment=PLAN1_DIST=/home/clownbot/srv/plan1
EnvironmentFile=/home/clownbot/plan1/.env
ExecStart=/usr/local/bin/deno run --no-lock --allow-read --allow-net --allow-env --allow-run --allow-write /home/clownbot/plan1/server.js
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
