# deploy/plan1.service
[Unit]
Description=plan1 server (clownbot)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/plan1
# resolve via login shell PATH, not a hardcoded install path — deno lives
# at ~/.deno/bin on some machines, /usr/local/bin on others (confirmed
# both exist in the wild across this project's own hosts)
ExecStart=/bin/bash -lc 'deno task serve'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
